# eve/workers/build/supervisor/livy.conf
[supervisord]

[program:livy]
command=/bin/bash /app/livy/entrypoint.sh
logfile=/etc/supervisor/conf.d/supervisor.log
autostart=true
autorestart=false
startretries=1
stdout_logfile=/etc/supervisor/conf.d/livy.log
stdout_logfile_maxbytes=0
stderr_logfile=/etc/supervisor/conf.d/livy.log
stderr_logfile_maxbytes=0
