# eve/workers/build/supervisor/livy.conf
[supervisord]

[program:livy]
command=/bin/bash /app/livy/entrypoint.sh
logfile=/etc/supervisor/conf.d/supervisor.log
autostart=true
autorestart=false
