# eve/workers/build/supervisor/livy.conf
[supervisord]

[program:livy]
command=/bin/sh/ /app/livy/entrypoint.sh
autostart=true
autorestart=false
